# deploy/pvpccheap.service
# Systemd user service per PVPCCheap amb Podman Rootless
#
# Aquest fitxer s'instal·la automàticament amb:
#   sudo ./deploy/deploy.sh --setup
#
# Ubicació: ~/.config/systemd/user/pvpccheap.service
#
# Comandes útils:
#   systemctl --user start pvpccheap
#   systemctl --user stop pvpccheap
#   systemctl --user restart pvpccheap
#   systemctl --user status pvpccheap
#   journalctl --user -u pvpccheap -f

[Unit]
Description=PVPCCheap Backend Service (Rootless Podman)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pvpccheap
ExecStart=/usr/bin/podman-compose --env-file .env.production -f podman-compose.prod.yml up
ExecStop=/usr/bin/podman-compose --env-file .env.production -f podman-compose.prod.yml down
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=60

# Variables d'entorn per podman rootless
Environment="XDG_RUNTIME_DIR=/run/user/%U"

[Install]
WantedBy=default.target
